# dist/nexus.service
[Unit]
Description=WorkFort Nexus Daemon

[Service]
Type=exec
ExecStart=%h/.local/bin/nexusd
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
AmbientCapabilities=CAP_NET_ADMIN

[Install]
WantedBy=default.target
